# ansible/petclinic.service.j2
[Unit]
Description=Java Service
After=network.target

[Service]
ExecStart=/usr/bin/java -jar {{ jar_path }}/{{ jar_file }}
User=petclinic
Restart=on-failure

[Install]
WantedBy=multi-user.target
